# debian/mouseemu/DEBIAN/postinst
#! /bin/sh
# postinst script for mouseemu
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package
#

install_udev() {
    # install the udev file only once:
    #    - the first time the package is installed
    # OR - the first time the package is upgraded from a version earlier than 1.0.15-3
    if [ -z "$2" ] || dpkg --compare-versions "$2" le 0.15-3; then
        ln -sf ../mouseemu.rules /etc/udev/rules.d/025_mouseemu.rules
	if [ -x /sbin/udevcontrol ] ; then
	    # udevcontrol fails if udev is not running, but installation should not abort
	    udevcontrol reload_rules || true
	fi
    fi
}

case "$1" in
    configure)
       if [ -z $2 ]; then
           # First install of mouseemu
           if [ ! -d /dev/.udev/ ]; then
               if [ -e /sbin/MAKEDEV ]; then
                   (cd /dev; MAKEDEV uinput)
               else
                   echo "MAKEDEV not installed, skipping device node creation."
               fi
           else
               echo "udev active, skipping device node creation."
	   fi
       fi

       #install the udev rules symlink
       install_udev

       ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

# Automatically added by dh_installinit/11.3.2
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if [ -x "/etc/init.d/mouseemu" ]; then
		update-rc.d mouseemu defaults >/dev/null
		invoke-rc.d mouseemu start || exit 1
	fi
fi
# End automatically added section


exit 0


